# docs/spikes/assets/session-1-ec2/agentic-api.service
# Personal spike only — not curriculum. See SPIKE-002.
# Install: sudo cp agentic-api.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now agentic-api

[Unit]
Description=Agentic Engineering Demo 1 API (personal spike)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/agentic/agentic-engineering-in-practice
EnvironmentFile=/opt/agentic/agentic-engineering-in-practice/.env
ExecStart=/home/ubuntu/.local/bin/uv run uvicorn app.main:app --app-dir src/backend --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
